# _archive/services/bazarr-stop.timer
[Unit]
Description=Stop bazarr
# BindsTo=bazarr.service

[Timer]
OnCalendar=*-*-* 22:10:00
# Persistent=true

[Install]
WantedBy=timers.target
